# scripts/com.ryan.fooocus-autocopy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ryan.fooocus-autocopy</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/ryantaylorvegh/.openclaw/workspace/scripts/fooocus_autocopy.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>/Users/ryantaylorvegh/.openclaw/workspace/memory/fooocus-autocopy.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/ryantaylorvegh/.openclaw/workspace/memory/fooocus-autocopy.stderr.log</string>
</dict>
</plist>
